# linux/lab02/crontab
0 */4 * * * /bin/sh /home/student/Documents/lab02/logcleaner.sh
